# 99-rapl-permissions.rules
SUBSYSTEM=="powercap", KERNEL=="intel-rapl:0", RUN+="/bin/chmod o+r /sys/class/powercap/intel-rapl:0/energy_uj"
